# bundles/linux/etc/udev/lwr-usb-media.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="/usr/local/bin/lwr-new-media"
ACTION=="change", SUBSYSTEM=="block", RUN+="/usr/local/bin/lwr-new-media"
