# bundles/linux/debian/11/config/includes.chroot/etc/systemd/system/autologin-setup.service
[Unit]
Before=getty@tty1.service
After=live-config.service
ConditionPathExists=!/etc/luwrain/default-user


[Service]
Type=oneshot
ExecStart=/usr/local/bin/autologin-setup



[Install]
WantedBy=multi-user.target